# deploy/kindle-dash-render.service
[Unit]
Description=Render Kindle dashboard pages
[Service]
Type=oneshot
ExecStart=/opt/kindle-dash/render.sh
